# ftp_server/supervisord_ftp.conf
[unix_http_server]
file=/var/run/supervisor.sock  

[inet_http_server]         
port=*:9001        
username=            
password=


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisord]
nodaemon=true
user=root
loglevel=DEBUG

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 



[program:ftpd]
command=/usr/sbin/proftpd -n
priority=12
redirect_stderr=true
stdout_logfile=/var/log/dockemu/%(program_name)s

[program:olsrd]
command= /usr/sbin/olsrd -i eth0 -d 1 
redirect_stderr=true
stdout_logfile=/var/log/dockemu/%(program_name)s